# services/goonclaw-automaton/deploy/gce/goonclaw-autonomous.service
[Unit]
Description=GoonClaw Autonomous Vertex Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/goonclaw
Environment=NODE_ENV=production
ExecStart=/usr/bin/env npx tsx services/goonclaw-automaton/runtime-loop.ts
Restart=always
RestartSec=10
User=goonclaw
Group=goonclaw

[Install]
WantedBy=multi-user.target
